# wasla/scripts/backup/systemd/wasla-backup-healthcheck.service
[Unit]
Description=Wassla backup health verification
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
EnvironmentFile=/etc/wasla/backup.env
ExecStart=/var/www/wasla/wasla/scripts/backup/backup_health_check.sh
